# salt/profile/postfix/files/etc/systemd/system/postfix-queue-size-metrics.timer.j2
{{ pillar['managed_by_salt'] }}

[Unit]
BindsTo=postfix.service

[Timer]
OnCalendar=*:0/5
AccuracySec=1us

[Install]
WantedBy=timers.target
