# salt/profile/calendar/files/calendar-sidekiq.service
[Unit]
Description=Calendar Sidekiq job queue
After=network.target

[Service]
Type=simple

User=calendar

WorkingDirectory=/srv/www/calendar-o-o

Environment="RAILS_ENV=production"
Environment=REDIS_URL=unix:///run/redis/calendar.sock

ExecStart=/usr/bin/bundler.{{ ruby }} exec sidekiq

Restart=always

[Install]
WantedBy=multi-user.target
