# salt/profile/gitlab_runner/files/etc/systemd/system/docker-cleanup.timer.jinja
{{ pillar['managed_by_salt'] }}

[Unit]
Description=Scheduled cleanup of Docker stuff

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
