# salt/profile/gitlab_runner/files/etc/systemd/system/docker-cleanup.service.jinja
{{ pillar['managed_by_salt'] }}

[Unit]
Description=Clean up Docker stuff

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/cleanup.sh
