# salt/profile/monitoring/prometheus/files/textfile/systemd/postfix-queue-size-metrics.service.jinja
{{ pillar['managed_by_salt'] }}

[Service]
ExecStart=/usr/local/libexec/systemd/postfix-queue-size-metrics.sh
LogLevelMax=notice
