# salt/profile/monitoring/prometheus/files/textfile/systemd/zypper-metrics.timer.jinja
{{ pillar['managed_by_salt'] }}

[Unit]
Description=Metrics scheduler for zypper

[Timer]
OnCalendar=hourly
AccuracySec=1us

[Install]
WantedBy=timers.target
