# salt/profile/monitoring/prometheus/files/textfile/systemd/zypper-metrics.service.jinja
{{ pillar['managed_by_salt'] }}

[Unit]
Description=Metrics collector for zypper

[Service]
ExecStart=/usr/local/libexec/systemd/zypper-metrics.sh
